# ops/systemd/kiro-sync-identity.timer
[Unit]
Description=Daily Identity Center sync at 02:15 UTC (= 11:15 KST)
# ingest 보다 15분 먼저 — sync 가 학생 행 채워야 그 학생 사용량을 ingest 가 올바른 학교로 매핑.

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
# 부팅 직후 못 받은 일정이면 한 번 따라잡음
Persistent=true
Unit=kiro-sync-identity.service

[Install]
WantedBy=timers.target
