# ops/systemd/kiro-sync-identity.timer
[Unit]
Description=Daily Identity Center sync at 15:00 UTC (= 00:00 KST)
# 자정에 학생 명단 먼저 정리. 12시간 뒤 ingest 가 새 학생도 매핑 가능.
# sync 자체는 빠르게 끝남 (IC API 한 번 + DB upsert).

[Timer]
OnCalendar=*-*-* 15:00:00 UTC
# 부팅 직후 못 받은 일정이면 한 번 따라잡음
Persistent=true
Unit=kiro-sync-identity.service

[Install]
WantedBy=timers.target
