# ops/systemd/kiro-ingest.timer
[Unit]
Description=Daily Kiro ingest at 02:30 UTC (= 11:30 KST)
# Kiro 가 02:00 UTC 에 CSV 떨어뜨림. 약간의 여유 두고 02:30 시작.
# sync-identity-center 가 먼저 02:15 에 끝나있어야 학생 매핑이 최신.

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
Persistent=true
Unit=kiro-ingest.service

[Install]
WantedBy=timers.target
