# ops/systemd/kiro-ingest.timer
[Unit]
Description=Daily Kiro ingest at 02:05 UTC (= 11:05 KST)
# Kiro 가 02:00 UTC (11:00 KST) 에 전날 CSV 떨어뜨림. 5분 여유 후 11:05 KST 에 가져옴.
# 스냅샷 재계산까지 보통 30초 안에 끝남 → 11:05 시작 → 11:06 즈음 완료.
# sync-identity-center 는 같은 날 자정 (KST) 에 이미 끝나있어 학생 매핑 최신 상태.

[Timer]
OnCalendar=*-*-* 02:05:00 UTC
Persistent=true
Unit=kiro-ingest.service

[Install]
WantedBy=timers.target
